# contrib/muninndb.service
# MuninnDB systemd service unit
#
# Install:
#   cp muninndb.service /etc/systemd/system/muninndb.service
#   # Edit the file below to set your paths and credentials, then:
#   systemctl daemon-reload
#   systemctl enable muninndb
#   systemctl start muninndb
#
# Assumptions:
#   - Binary built from source at /opt/muninndb-src/muninndb-server
#   - Data directory at /opt/muninndb/data
#   - Adjust paths and Environment= values to match your setup
#
# After updating credentials, reload with:
#   systemctl daemon-reload && systemctl restart muninndb
#
# WARNING: lifecycle conflict — do NOT run 'muninn start', 'muninn stop',
# or 'muninn restart' while this service is enabled. These CLI commands and
# systemd both own the process; running them concurrently races on the
# Pebble database lock and can put systemd into a permanent Restart=on-failure
# loop. While this unit is active, use systemctl to control MuninnDB:
#   systemctl start muninndb   (instead of: muninn start)
#   systemctl stop muninndb    (instead of: muninn stop)
#   systemctl restart muninndb (instead of: muninn restart)
#   systemctl status muninndb  (instead of: muninn status)

[Unit]
Description=MuninnDB Memory Server
After=network.target

[Service]
Type=forking
# PIDFile must match your MUNINNDB_DATA path — muninn writes muninn.pid
# inside the data directory on start and removes it on clean stop.
PIDFile=/opt/muninndb/data/muninn.pid

# --- Data ---
# MUNINNDB_DATA overrides the default (~/.muninn/data).
Environment=MUNINNDB_DATA=/opt/muninndb/data

# --- Network ---
Environment=MUNINN_LISTEN_HOST=0.0.0.0
Environment=MUNINN_CORS_ORIGINS=http://YOUR_SERVER_IP:8476

# --- Embeddings (choose one, comment out the rest) ---
# Local Ollama (no API key required):
# Environment=MUNINN_OLLAMA_URL=ollama://localhost:11434/nomic-embed-text
# OpenAI:
Environment=MUNINN_OPENAI_KEY=sk-...
# Voyage AI:
# Environment=MUNINN_VOYAGE_KEY=...

# --- LLM Enrichment (optional) ---
# Anthropic:
Environment=MUNINN_ENRICH_URL=anthropic://claude-haiku-4-5-20251001
Environment=MUNINN_ANTHROPIC_KEY=sk-ant-...
# OpenAI:
# Environment=MUNINN_ENRICH_URL=openai://gpt-4o-mini
# Environment=MUNINN_ENRICH_API_KEY=sk-...

# --- TLS (optional) ---
# Serve clients over https. Set BOTH or neither. See docs/tls.md.
# NOTE: ~/.muninn/muninn.env is NOT loaded under systemd (no $HOME), so TLS
# for a systemd deployment MUST be set here via Environment= (or EnvironmentFile=).
# Environment=MUNINN_TLS_CERT=/etc/muninndb/cert.pem
# Environment=MUNINN_TLS_KEY=/etc/muninndb/key.pem
# Alternatively load from a separate file (systemd EnvironmentFile format is
# KEY=VALUE, NOT muninn.env format): EnvironmentFile=-/etc/muninndb/muninndb.env

# --- Auth ---
# The MCP bearer token is passed via --mcp-token flag.
# Alternatively, write your token to ~/.muninn/mcp.token (no flag needed):
#   echo 'your-token' > ~/.muninn/mcp.token && chmod 600 ~/.muninn/mcp.token
ExecStart=/opt/muninndb-src/muninndb-server start --mcp-token your-mcp-token-here
ExecStop=/opt/muninndb-src/muninndb-server stop
Restart=on-failure

[Install]
WantedBy=multi-user.target
